# installers/macos/com.aura.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aura.api</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/share/dotnet/dotnet</string>
        <string>run</string>
        <string>--project</string>
        <string>/usr/local/share/aura/src/Aura.Api</string>
        <string>--no-launch-profile</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/share/aura</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>ASPNETCORE_ENVIRONMENT</key>
        <string>Production</string>
        <key>ASPNETCORE_URLS</key>
        <string>http://localhost:5300</string>
        <key>ConnectionStrings__auradb</key>
        <string>Host=localhost;Port=5432;Database=aura;Username=aura;Password=aura</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/aura/api.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/aura/api.error.log</string>
</dict>
</plist>
